# agent/templates/com.qwickapps.qwickguard-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.qwickapps.qwickguard-agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>__VENV_PATH__/bin/qwickguard-agent</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>QWICKGUARD_CONFIG</key>
        <string>__CONFIG_PATH__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/.qwickguard/logs/agent.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.qwickguard/logs/agent.stderr.log</string>
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
